# crontab
#MAILTO=email@example.com

0 15 * * * * /bin/bash /script.sh
